# build/overlays-t1/system/etc/systemd/system/usb-mount@.service
# FLSUN T1 Open Source Edition

[Unit]
Description=Mount USB Drive on %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/home/pi/flsun-os/system/usb-mount.sh add %i
ExecStop=/home/pi/flsun-os/system/usb-mount.sh remove %i
